# bot/katana-executor.service
# systemd unit for the live executor. Resilient autostart + anti-crashloop.
# Install (needs sudo):
#   sudo cp bot/katana-executor.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now katana-executor
# Logs:  journalctl -u katana-executor -f
# Stop:  sudo systemctl stop katana-executor
#
# Live mode is enabled in ~/.katana-bot/env: DRY_RUN=0 (code default is 1 = safe).
# The kill-switch stops the loop (returns) on a limit breach — then investigate,
# `python3 -m bot.executor reset`, and `systemctl start katana-executor`.

[Unit]
Description=Katana Morpho liquidator executor (chainId 747474)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude-agent
WorkingDirectory=/home/claude-agent/katana-liquidator
EnvironmentFile=/home/claude-agent/.katana-bot/env
ExecStart=/usr/bin/python3 -u -m bot.executor loop
Restart=on-failure
RestartSec=10
# anti-crashloop: at most 5 restarts / 5 min, else the unit gives up (needs intervention)
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
